# systemd/zeta-singbox.service
[Unit]
Description=ZetaVPN sing-box core
Documentation=https://sing-box.sagernet.org
After=network-online.target zeta-panel.service
Wants=network-online.target
# Cap a permanent crash-loop (see zeta-xray.service): failed after 20 restarts
# in 5 min instead of looping forever; transient failures still auto-recover.
StartLimitIntervalSec=300
StartLimitBurst=20

[Service]
Type=simple
User=zetavpn
Group=zetavpn
WorkingDirectory=/etc/sing-box
ExecStart=/usr/local/bin/sing-box run -c /etc/sing-box/config.json
Restart=on-failure
RestartSec=3
LimitNOFILE=1000000
# install.sh chowns /etc/sing-box and the cert dir to 'zetavpn' so this
# non-root user can write its cache file and read its config + TLS certs.
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_NET_ADMIN
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
# sing-box writes its Clash-API cache file under here (experimental.cache_file
# in the generated config) — the only path it needs write access to.
ReadWritePaths=/etc/sing-box

[Install]
WantedBy=multi-user.target
