# systemd/zeta-ws.service
[Unit]
Description=ZetaVPN SSH-over-WebSocket proxy
After=network-online.target ssh.service
Wants=network-online.target
# Cap a permanent crash-loop (see zeta-xray.service): failed after 20 restarts
# in 5 min instead of looping forever; transient failures still auto-recover.
StartLimitIntervalSec=300
StartLimitBurst=20

[Service]
Type=simple
User=zetavpn
Group=zetavpn
ExecStart=/usr/bin/python3 /opt/zetavpn/proxies/ws-proxy.py --listen 0.0.0.0:8880 --backend 127.0.0.1:22
Restart=always
RestartSec=3
LimitNOFILE=1000000
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
