# systemd/zeta-panel.service
[Unit]
Description=ZetaVPN control panel (FastAPI)
After=network-online.target
Wants=network-online.target
# Cap a permanent crash-loop (see zeta-xray.service): failed after 20 restarts
# in 5 min instead of looping forever; transient failures still auto-recover.
StartLimitIntervalSec=300
StartLimitBurst=20

[Service]
Type=simple
User=zetavpn
Group=zetavpn
WorkingDirectory=/opt/zetavpn/backend
EnvironmentFile=-/opt/zetavpn/.env
ExecStart=/opt/zetavpn/venv/bin/python -m zeta.main
Restart=always
RestartSec=3
LimitNOFILE=1000000
# systemd creates + owns /var/lib/zeta-ssh-info (0755, zetavpn:zetavpn) on every
# start, so the panel can write each SSH account's post-auth banner file there and
# the unprivileged tunnel login shell — which can't traverse the 0750 ZETA_HOME —
# can read its own. Absolute path is /var/lib/<name>; matches config.ssh_info_dir.
StateDirectory=zeta-ssh-info
# The panel runs as the unprivileged 'zetavpn' user (install.sh creates it and
# chowns the dirs it writes: ZETA_HOME, /usr/local/etc/xray, /etc/sing-box).
# The few genuinely-privileged jobs it still needs — managing SSH tunnel
# accounts (useradd/chpasswd/chage/usermod/userdel) and reloading the proxy
# core / SSH-stack services (systemctl) — are delegated via a narrowly scoped
# NOPASSWD sudoers rule (/etc/sudoers.d/zetavpn-panel) instead of running the
# whole HTTP-facing app as root. NoNewPrivileges MUST stay false: sudo is a
# setuid-root binary, and NoNewPrivileges=true would block it from working.
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
